# sv/goagent/run
#!/bin/sh

sv start goserver || exit 1

exec 2>&1
exec /usr/share/go-agent/agent.sh
